# systemd/01.ondemand/example.socket
[Unit]
Description=Example service socket

[Socket]
ListenStream=/tmp/example.sock

[Install]
WantedBy=sockets.target
