# systemd/03.onboot/example.service
[Unit]
Description=Example service

[Service]
ExecStart=/usr/bin/example

[Install]
Also=example.socket
WantedBy=multi-user.target
