# systemd/01.example/example.service
[Unit]
Description=Example service

[Service]
ExecStart=/bin/sleep 100

[Install]
Also=example.socket
#WantedBy=multi-user.target
